# systemd/ray-head.service
[Unit]
Description=Ray Head Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root

Environment=RAY_USAGE_STATS_ENABLED=0
Environment=RAY_TMPDIR=/tmp/ray
Environment=TMPDIR=/tmp/ray
Environment=TEMP=/tmp/ray
Environment=TMP=/tmp/ray

LimitNOFILE=1048576
LimitNPROC=65536

ExecStart=/opt/ray/bin/ray start --head --block \
  --port=6379 \
  --dashboard-host=0.0.0.0 \
  --dashboard-port=8265 \
  --object-manager-port=8076 \
  --object-spilling-directory=/data/ray/spill \
  --object-store-memory=2147483648 \
  --num-cpus=8 \
  --resources='{"head": 1}' \
  --metrics-export-port=8080

ExecStop=/opt/ray/bin/ray stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
